# crontab
0-59/15 * * * *        for repo in ~/git/*; do git -C $repo fetch origin > /dev/null 2>&1; done
